# util/templates/env.service
[Unit]
Description={{ description }}

[Service]
Type=oneshot
ExecStart={{ script }} {{ ip }} {{ location }} --log {{ log }}
