# util/templates/oncalender.timer
[Unit]
Description={{ description }}

[Timer]
OnCalendar={{ schedule }}
Persistent={{ persistence }}

[Install]
WantedBy=timers.target
