# services/env-node/env-logger.timer
[Unit]
Description=Timer for env-logger.unit

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
