# services/env-node/env-logger.timer
[Unit]
Description=Timer for env-logger.unit

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
